# src/diepxuan/management/ductn/etc/cron.daily/ductn
#!/usr/bin/env bash
#!/bin/bash

ductn update 2>/dev/null
