# src/diepxuan/management/etc/cron.daily/ductn
#!/usr/bin/env bash
#!/bin/bash

ductn selfupdate 2>/dev/null
